# .deploy_scripts(DO NOT USE)/clear_sessions.timer
[Unit]
Description=Timer for crear user's sessions

[Timer]
OnBootSec=300
OnCalendar=monthly
Persistent=True

[Install]
WantedBy=multi-user.target
