# .deploy_scripts(DO NOT USE)/star_burger_docker_db.service
[Service]
Type=oneshot
ExecStart=/usr/bin/docker start star_burger_db

[Install]
WantedBy=multi-user.target
